# lab7/roles/bind/templates/bind_exporter.service.j2
[Unit]
Description=Bind Exporter

[Service]
User=root
ExecStart=/usr/local/bin/bind_exporter

[Install]
WantedBy=multi-user.target
