# exam/roles/mysql_exporter/templates/mysql_exporter.service.j2
[Unit]
Description=Mysql exporter

[Service]
User=root
ExecStart=/opt/mysql_exporter/mysqld_exporter-0.12.1.linux-amd64/mysqld_exporter

[Install]
WantedBy=multi-user.target
